# roles/archlinux/files/pacman-sync.service
[Unit]
Description=update pacman databases

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Syy
