# roles/archlinux/files/pacman-sync.timer
[Unit]
Description=update pacman databases on boot

[Timer]
OnBootSec=15min

[Install]
WantedBy=timers.target
