# roles/myscripts/templates/nyx.kumobackup.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>nyx.kumobackup</string>
    <key>WorkingDirectory</key><string>{{ ansible_env.HOME }}/Library/Logs/kumobackup</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>fr_FR.UTF-8</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>{{ ansible_env.HOME }}/bin/kumobackup</string>
        <string>backup</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
        <key>Weekday</key>
        <integer>5</integer>
    </dict>
    <key>StandardErrorPath</key> <string>stderr.log</string>
    <key>StandardOutPath</key>   <string>stdout.log</string>
</dict>
</plist>